# scripts/com.ryan.drain-queue.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Approval Queue consumer (RYA-166).

  Every 15 minutes: scan `2. Areas/Social Media/Approval Queue/` and schedule
  anything Ryan marked `status: approved` via the Blotato REST API. The vault is in
  iCloud and syncs to Obsidian on his phone, so this timer IS the phone-approval
  path: edit one word on the couch, it ships within 15 minutes. No app required —
  which is why this is launchd and not a watcher inside RyanOS.

  NOTE the --ship flag. Without it the script is a dry run; this is the one caller
  allowed to actually schedule. It still cannot approve anything.

  Runs via drain-queue.sh, not bun directly, because the wrapper resolves
  BLOTATO_API_KEY from ~/.secrets.zsh. `op run` is NOT used here: it needs an
  interactive unlock, and launchd has no one to touch the sensor.

  PRECONDITION — this timer no-ops until the key is resolvable:
    echo 'export BLOTATO_API_KEY="..."' >> ~/.secrets.zsh    # see drain-queue.sh

  Install:
    cp scripts/com.ryan.drain-queue.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.ryan.drain-queue.plist
  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.ryan.drain-queue.plist
  Watch:
    tail -f ~/.ade/drain-queue.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ryan.drain-queue</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ryannaquin/Code/damon-ade/scripts/drain-queue.sh</string>
        <string>--ship</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ryannaquin/Code/damon-ade</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/ryannaquin/.ade/drain-queue.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ryannaquin/.ade/drain-queue.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ryannaquin/.bun/bin:/Users/ryannaquin/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
